# debian/untangle-classd.prerm
#!/bin/bash

LOG_PATH=/var/log/untangle-classd

/etc/init.d/untangle-classd stop
update-rc.d untangle-classd remove > /dev/null
rm -f $LOGPATH/*

exit 0

