# debian/untangle-classd.postinst
#!/bin/bash

LOG_FILE=/var/log/untangle-classd/classd.log

-mkdir -p `dirname ${LOG_FILE}`

echo "[`date`] DEBUG untangle-classd installed or upgraded" >> ${LOG_FILE}

if [ -x "/etc/init.d/untangle-classd" ]; then
    update-rc.d untangle-classd defaults 25 > /dev/null
    /etc/init.d/untangle-classd restart
fi

exit 0
